# vagrant/services/h2-http.service
[Unit]
Description=h2 HTTP service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/h2-http/server.py
Restart=on-failure
RestartSec=2
User=h2-http
Group=h2-http
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
